# scripts/services/medical-deid-backend.service
[Unit]
Description=Medical De-identification Backend API
After=network.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=__APP_USER__
Group=__APP_GROUP__
WorkingDirectory=__PROJECT_ROOT__/web/backend
EnvironmentFile=/etc/medical-deid/medical-deid.env
Environment="PATH=__PROJECT_ROOT__/.venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=__PROJECT_ROOT__/.venv/bin/uvicorn main:app --host ${MEDICAL_DEID_BACKEND_HOST} --port 8000
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/lib/medical-deid /var/log/medical-deid __PROJECT_ROOT__/web/backend/data __PROJECT_ROOT__/web/backend/logs __PROJECT_ROOT__/logs
CapabilityBoundingSet=
LockPersonality=true
RestrictSUIDSGID=true
UMask=0077

[Install]
WantedBy=multi-user.target
