# scripts/services/medical-deid-frontend.service
[Unit]
Description=Medical De-identification Frontend (Production)
After=network.target medical-deid-backend.service
Wants=medical-deid-backend.service

[Service]
Type=simple
User=__SERVICE_USER__
Group=__SERVICE_GROUP__
WorkingDirectory=__PROJECT_ROOT__/web/frontend
Environment="PATH=__NVM_NODE_PATH__:/usr/local/bin:/usr/bin:/bin"
ExecStart=__NVM_NODE_PATH__/serve -s dist -l 127.0.0.1:5173
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
